# config/supervisord/rancher-gen.conf
[program:rancher-gen]
command = /usr/local/bin/rancher-gen --config=/etc/rancher-gen.cfg
stdout_events_enabled=true
stderr_events_enabled=true
autostart=true
autorestart=true
