# config/supervisord/rsyslog.ini
[program:rsyslogd]
command = /usr/sbin/rsyslogd -n
autostart=true
autorestart=true
